# autostart/.config/autostart/backup.desktop
[Desktop Entry]
Name=backup
GenericName=backup
Exec=/home/vago/bin/backup
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
