# autostart/.config/autostart/bup.desktop
[Desktop Entry]
Name=Backup
Type=Application
Exec=bup 
NoDisplay=true
